# pyro_1/supervisord.conf
[supervisord]
nodaemon=true

[program:pyro4]
command=pyro4-ns

[program:service_creator]
command=python3 /pyro-algorithms-orax/service_creator.py
